# crontab.cron
*/01 * * * * ~/2015-02-web-historian/workers/htmlfetcher.js
